# lumagen_ui/lumagen-ui.service
[Unit]
Description=Lumagen + JVC Web Control UI
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment="LUMAGEN_PORT=__LUMAGEN_PORT__"
Environment="JVC_HOST=__JVC_HOST__"
ExecStart=__VENV__/bin/gunicorn \
    --workers 1 \
    --bind 0.0.0.0:5001 \
    --timeout 30 \
    --access-logfile - \
    app:app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
